# crates/relay-server/cascade-relay.service
# systemd unit for cascade-relay.
#
# Install the unit file and the launcher script, then reload and enable:
#
#   install -m 0755 crates/relay-server/cascade-relay-start.sh /usr/local/bin/
#   cp crates/relay-server/cascade-relay.service /etc/systemd/system/
#   systemctl daemon-reload
#   systemctl enable --now cascade-relay
#
# Configuration lives in /etc/cascade/relay.env (create it before starting).
# At minimum it must set CASCADE_RELAY_SHARED_SECRET.  See that file's template
# at crates/relay-server/relay.env.example in the source tree, or read the
# inline comments on EnvironmentFile below.
#
# The ExecStart launcher script (cascade-relay-start.sh) conditionally appends
# --metrics-bind and --announce-bind only when the corresponding variables are
# non-empty.  systemd does not support ${VAR:-default} expansion in ExecStart,
# so default values and optional flags are handled in the script rather than
# inline.

[Unit]
Description=Cascade relay server
Documentation=https://github.com/Mearman/cascade/blob/main/docs/deployment.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# The binary reads the shared secret from the environment to avoid it
# appearing in /proc/<pid>/cmdline or process listings.
EnvironmentFile=/etc/cascade/relay.env

ExecStart=/usr/local/bin/cascade-relay-start.sh

# Restart on unexpected exit. 5-second back-off avoids tight loops if the
# port is held by a lingering previous instance.
Restart=on-failure
RestartSec=5s

# Drop privileges after startup. The relay needs no special capabilities:
# it only listens on port 9999 (> 1024) and opens outbound connections.
# If you want to listen below port 1024, add AmbientCapabilities=CAP_NET_BIND_SERVICE
# and remove the User/Group lines, or use a reverse proxy instead.
User=cascade-relay
Group=cascade-relay
NoNewPrivileges=true

# Reduce the writable filesystem surface. The relay holds no persistent
# state (no files, no database), so these restrictions are safe.
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true

# Allow binding to the TCP listener. The relay does not need other network
# namespaces or privileges.
RestrictAddressFamilies=AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
